# tests/integration/functional/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:rustybgpd]
command=/usr/bin/rustybgpd
